# webService/seafile.service
[Unit]
Description=seafile file storge auto launcher

[Service]
Type=forking
ExecStart=/cloudStorge/launch.sh

[Install]
WantedBy=multi-user.target
